# scripts/trading-agent-dashboard.service
[Unit]
Description=Trading-Agent Streamlit Dashboard
After=network.target

[Service]
Type=simple
User=hermes
WorkingDirectory=/home/hermes/workspace/trading-agent
ExecStart=/home/hermes/.local/bin/uv run streamlit run tradingagents/dashboard/app.py --server.headless=true --server.port=8501 --server.address=0.0.0.0
Restart=on-failure
RestartSec=5
Environment=PATH=/home/hermes/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
